# services/php-fpm/run
#!/bin/bash

exec /usr/sbin/php-fpm5.6 -F -y /config/nginx/nginx-fpm.conf
